# packages/fasva-osmserver/fasva_osmserver-0.1.6-py3-none-any.whl/fasva_osmserver-0.1.6.data/data/share/osmserver.service
[Unit]
Description=FASva-OSMServer startup script

[Service]
ExecStart=<root>osmserver --config=share/example.cfg --log-level=INFO
Type=simple

[Install]
WantedBy=graphical.target
